# ctf-discord-bot.service
[Unit]
Description = ctf-discord-bot

[Service]
ExecStart = /root/ctf-discord-bot/app.py
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target