# ctf-discord-bot.service
[Unit]
Description = ctf-discord-bot

[Service]
ExecStart = <YOUR_HOME_DIR_PATH>/ctf-discord-bot/app.py
Restart = always
Type = simple
Environment = "DISCORD_TOKEN=<YOUR_TOKEN_HERE>"

[Install]
WantedBy = multi-user.target
