# ctf-discord-bot.service
[Unit]
Description = ctf-discord-bot

[Service]
Type = simple
WorkingDirectory=<YOUR_HOME_DIR_PATH>/ctf-discord-bot
ExecStart = <YOUR_HOME_DIR_PATH>/.local/bin/uv run main.py
Restart = always
Environment = "DISCORD_TOKEN=<YOUR_TOKEN_HERE>"

[Install]
WantedBy = multi-user.target
